# triggers/launchd/com.shupei.agent.telegram_bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Long-polling Telegram bot. KeepAlive auto-restarts on crash;
         ThrottleInterval keeps a flapping process from spamming the Telegram
         API. No FDA required. -->
    <key>Label</key>
    <string>com.shupei.agent.telegram_bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-u</string>
        <string>-m</string>
        <string>triggers.telegram_bot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/smo/repos/home-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/smo</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/smo/repos/home-agent/data/runs/telegram_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/smo/repos/home-agent/data/runs/telegram_bot.err</string>
</dict>
</plist>
